# eszett.service
[Unit]
Description=eszett

[Service]
Type=simple
WorkingDirectory=/TBD
ExecStart=/usr/bin/pipenv run eszett
KillSignal=SIGINT
